# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/app/supervisor/supervisord.pid
loglevel=info

[program:discovery-server]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/DiscoveryServer-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=30
stderr_logfile=/app/logs/discovery-server.err.log
stdout_logfile=/app/logs/discovery-server.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=100

[program:config-server]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/ConfigServer-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=20
stderr_logfile=/app/logs/config-server.err.log
stdout_logfile=/app/logs/config-server.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=200

[program:api-gateway]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/ApiGateway-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=20
stderr_logfile=/app/logs/api-gateway.err.log
stdout_logfile=/app/logs/api-gateway.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=300

[program:user-service]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/UserService-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=20
stderr_logfile=/app/logs/user-service.err.log
stdout_logfile=/app/logs/user-service.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=300
environment=DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s"

[program:recharge-service]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/RechargeService-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=20
stderr_logfile=/app/logs/recharge-service.err.log
stdout_logfile=/app/logs/recharge-service.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=300
environment=DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s"

[program:payment-service]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/PaymentService-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=20
stderr_logfile=/app/logs/payment-service.err.log
stdout_logfile=/app/logs/payment-service.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=300
environment=DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s"

[program:notification-service]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/NotificationService-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=20
stderr_logfile=/app/logs/notification-service.err.log
stdout_logfile=/app/logs/notification-service.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=300
environment=DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s"

[program:operator-service]
command=java %(ENV_JAVA_OPTS)s -jar /app/services/OperatorService-1.0.jar
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=20
stderr_logfile=/app/logs/operator-service.err.log
stdout_logfile=/app/logs/operator-service.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
priority=300
environment=DB_HOST="%(ENV_DB_HOST)s",DB_PORT="%(ENV_DB_PORT)s",DB_NAME="%(ENV_DB_NAME)s",DB_USER="%(ENV_DB_USER)s",DB_PASSWORD="%(ENV_DB_PASSWORD)s"

[unix_http_server]
file=/app/supervisor/supervisor.sock

[supervisorctl]
serverurl=unix:///app/supervisor/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
